# unegui/crontab
50 23 * * * cd /app && python ./main.py >> /logfile